# etc/profile.d/home.sh
[ -d $HOME/bin ] && PATH=$PATH:$HOME/bin
export PATH

